# ova_hub.service.sample
[Unit]
Description=Run Open Voice Assistant Hub
After=multi-user.target

[Service]
ExecStart=/home/pi/openvoiceassistant_hub/start_hub.sh
User=user

[Install]
WantedBy=multi-user.target