# scripts/ova_hub_backend.service.sample
[Unit]
Description=openvoiceassistant hub frontend
After=multi-user.target

[Service]
ExecStart={$pwd}/scripts/start_backend.sh
User=user

[Install]
WantedBy=multi-user.target